# stubs/docker/servd/build/config/supervisord.conf
[supervisord]
nodaemon=true
user=servd
logfile=/dev/null
logfile_maxbytes=0
pidfile=/home/servd/supervisord.pid

[program:php-fpm]
command=php-fpm -F
user=servd
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
user=servd
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
